# tasks/avail.cron
#!/bin/bash

export PATH=/usr/local/bin:$PATH

cd /opt/mobile/tasks && source ./setting && ./avail.task

